# example/environment/crontab
@reboot searchd
@reboot indexer --all --rotate

* * * * * indexer magnet --rotate